# roles/nextcloud/templates/systemd/certbot_renewal.service.j2
[Unit]
Description=Certbot certificate renewal

[Service]
ExecStart=/usr/bin/certbot renew --pre-hook "service nginx stop" --post-hook "service nginx start"

[Install]
WantedBy=basic.target
